# systemd-user/.config/systemd/user/app-com.mitchellh.ghostty.service.d/override.conf
[Service]
# Ghostty can exit 0 before sd_notify handshake completes in some builds.
# Type=simple avoids false protocol failures while keeping preload behavior.
Type=simple
BusName=
